# selfhost/memos/user_units/memos-proxy.socket
[Unit]
Description=Socket for on-demand memos proxy (Port 5230).

[Socket]
ListenStream=5230
Accept=no

[Install]
WantedBy=sockets.target