# selfhost/languagetool/user_units/languagetool-proxy.socket
[Unit]
Description=Socket for on-demand languagetool proxy (Port 8081).

[Socket]
ListenStream=8081
Accept=no

[Install]
WantedBy=sockets.target